# scripts/macos/launchd/com.kiosky.player.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.kiosky.player</string>
    <key>ProgramArguments</key>
    <array>
      <string>/path/to/kioskyLocal/.venv/bin/python</string>
      <string>/path/to/kioskyLocal/kiosk.py</string>
      <string>--config</string>
      <string>/path/to/kioskyLocal/config.json</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/kioskyLocal</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/path/to/kioskyLocal/logs/kiosky.out.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/kioskyLocal/logs/kiosky.err.log</string>
  </dict>
</plist>
